# roles/motd/templates/motd-update-cache.service.j2
# -----------------------------------------------------------------------------
# THIS FILE IS MANAGED BY ANSIBLE — DO NOT EDIT MANUALLY
#
# Purpose:
#   Oneshot service to refresh the MOTD package update cache.
#   Hardened with systemd security directives.
#
# Source:
#   Template: {{ template_path }}
#   Control:  {{ template_host }} ({{ template_host }}.its.ax)
#
# Changes:
#   Edit the template and re-run Ansible.
#   Manual changes will be overwritten on the next playbook run.
#
# Contact:
#   IT-stödperson <itstodperson@its.ax>
#
# Deployed to: {{ inventory_hostname }}
# Generated:   {{ template_run_date }}
# -----------------------------------------------------------------------------

[Unit]
Description=Refresh MOTD package update cache
After=network-online.target apt-daily.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-motd-cache.sh

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
PrivateDevices=yes
PrivateUsers=no
SystemCallArchitectures=native

# Allow writing to /run for the cache file
ReadWritePaths=/run

# Only allow necessary system calls
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount

# Resource limits
TimeoutStartSec=60
